# packages/sysutils/amremote/system.d/amremote.service
[Unit]
Description=Enable amremote and disable lirc

[Service]
Type=oneshot
ExecStart=/usr/lib/libreelec/meson_remote_ir meson-remote

[Install]
WantedBy=multi-user.target
